# persistence/com.globe-cli.server.plist
<!-- macOS launchd plist — ~/Library/LaunchAgents/com.globe-cli.server.plist -->
<!-- Install: cp com.globe-cli.server.plist ~/Library/LaunchAgents/ -->
<!--          launchctl load ~/Library/LaunchAgents/com.globe-cli.server.plist -->
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.globe-cli.server</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/globe-cli/.venv/bin/python</string>
        <string>/opt/globe-cli/server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/opt/globe-cli</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/opt/globe-cli/logs/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/globe-cli/logs/stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
